# tools/unseal-page/openclaw-unseal-page.service
[Unit]
Description=OpenClaw unseal page
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/sai/.openclaw/workspace/tools/unseal-page
ExecStart=/usr/bin/node /home/sai/.openclaw/workspace/tools/unseal-page/server.js
Restart=on-failure
RestartSec=2
Environment=UNSEAL_HOST=192.168.32.254
Environment=UNSEAL_PORT=18888
Environment=UNSEAL_STATE_DIR=/home/sai/.openclaw/workspace/state/unseal

[Install]
WantedBy=default.target
